# etc/profile.d/alces-adminware.sh
if id |grep -qw siteadmin ; then
  admin(){
    sudo /opt/adminware/adminware.rb
  }
fi
